# scripts/3lora-ltx23.service
[Unit]
Description=3lora LTX2.3 AutoDL Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/3lora-ltx23-app
ExecStart=/root/3lora-ltx23-app/scripts/improved-autostart.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# 环境变量
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
